# deploy/com.missioncontrol.security-chief-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.missioncontrol.security-chief-daily</string>
  <key>ProgramArguments</key><array>
    <string>/usr/local/bin/node</string>
    <string>/Users/terminal/.openclaw/workspace/mission-control/scripts/security-chief-daily.mjs</string>
  </array>
  <key>StartCalendarInterval</key><dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
  <key>RunAtLoad</key><false/>
  <key>StandardOutPath</key><string>/Users/terminal/.openclaw/workspace/mission-control/logs/security-chief-daily.out.log</string>
  <key>StandardErrorPath</key><string>/Users/terminal/.openclaw/workspace/mission-control/logs/security-chief-daily.err.log</string>
</dict></plist>
