# deploy/com.missioncontrol.twitter-orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.missioncontrol.twitter-orchestrator</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>/Users/terminal/.openclaw/workspace/mission-control/scripts/run-twitter-orchestrator.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/terminal/.openclaw/workspace/mission-control</string>
    <key>StandardOutPath</key>
    <string>/Users/terminal/.openclaw/workspace/mission-control/logs/twitter-orchestrator.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/terminal/.openclaw/workspace/mission-control/logs/twitter-orchestrator.err.log</string>
  </dict>
</plist>
